# systemd/traffic-local-bot.service
[Unit]
Description=traffic-local-notify telegram command listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/traffic-local/bot_listener.py
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
